# crontab
#  m h  dom mon dow    command
   * *  *   *   *      /repos/sync.sh
